# root/usr/lib/systemd/system/render-sops-secrets.service
[Unit]
Description=Render SOPS secrets from TPM-wrapped age identity
Documentation=man:systemd-creds(1) man:systemd.exec(5)

DefaultDependencies=no
Requires=enroll-sops-age-key.service
After=enroll-sops-age-key.service local-fs.target systemd-ask-password-console.service
Before=sysinit.target shutdown.target
Wants=systemd-ask-password-console.service
Conflicts=shutdown.target

ConditionPathExists=/etc/credstore.encrypted/sops-age-key.cred

[Service]
Type=oneshot

LoadCredentialEncrypted=sops-age-key:/etc/credstore.encrypted/sops-age-key.cred
Environment=SOPS_AGE_KEY_FILE=%d/sops-age-key
ExecStart=/usr/libexec/render-sops-secrets.sh

UMask=0077
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadOnlyPaths=/usr
ReadWritePaths=/etc
ReadWritePaths=/var
NoNewPrivileges=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

[Install]
WantedBy=multi-user.target
