# root/etc/profile.d/tmout.sh
# skip check
# shellcheck shell=sh
declare -xr TMOUT=600
# vim: set ft=bash et tw=4 sw=4 sts=4:
