# root/etc/profile.d/usr-local.sh
# shellcheck shell=sh
pathmunge /usr/local/bin
